# .config/systemd/user/battery-notify.service
[Unit]
Description=Battery low notification

[Service]
ExecStart=%h/.config/hypr/scripts/battery_low.sh
